# scripts/raylu-slackbot.service
[Unit]
Description=raylu-tv Slack bot (Socket Mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /usr/local/bin/raylu-slackbot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Slice=sidejobs.slice

[Install]
WantedBy=multi-user.target
